# deploy/quadlet/mneme.container
[Unit]
Description=Mneme all-in-one (AI API with embedded SQLite + QQ bridge)

[Container]
Image=ghcr.io/vesperglow/mneme:latest
ContainerName=mneme
EnvironmentFile=./mneme.env
PublishPort=127.0.0.1:8000:8000
PublishPort=127.0.0.1:9000:9000
# 从旧的 qq-agent-* 卷迁移数据后再启用（见 README「切到 mneme」一节），
# 否则会挂到新建的空卷、丢掉既有 memory.db。
Volume=mneme-data:/data
# 纯文本进出，稳态只有几十 MB（axum + SQLite，无本地模型、无图片解码）。
PodmanArgs=--memory=256m
AutoUpdate=registry
HealthCmd=curl -fsS http://127.0.0.1:8000/health/live >/dev/null && curl -fsS http://127.0.0.1:9000/healthz >/dev/null
HealthInterval=15s
HealthRetries=5
# 无模型下载，启动就是开监听。
HealthStartPeriod=30s
HealthTimeout=5s
HealthOnFailure=kill
Notify=healthy

[Service]
Restart=always
RestartSec=5s
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=default.target
