# deploy/quadlet/qq-agent.container
[Unit]
Description=QQ Agent all-in-one (AI API with embedded SQLite/ONNX + BotGo bridge)

[Container]
Image=ghcr.io/vesperglow/qq-agent:latest
ContainerName=qq-agent
EnvironmentFile=./qq-agent.env
PublishPort=127.0.0.1:8000:8000
PublishPort=127.0.0.1:9000:9000
Volume=qq-agent-data:/data
Volume=qq-agent-models:/models
# embedding 权重约 650MB 常驻，稳态总占用 ~1GB；上限兜底 2GB。
PodmanArgs=--memory=2g
AutoUpdate=registry
HealthCmd=python -c "import urllib.request; urllib.request.urlopen('http://127.0.0.1:8000/health/live', timeout=3); urllib.request.urlopen('http://127.0.0.1:9000/healthz', timeout=3)"
HealthInterval=15s
HealthRetries=5
# 首次启动要下载 embedding 模型（约 640MB），给足冷启动时间。
HealthStartPeriod=300s
HealthTimeout=5s
HealthOnFailure=kill
Notify=healthy

[Service]
Restart=always
RestartSec=5s
TimeoutStartSec=600
TimeoutStopSec=30

[Install]
WantedBy=default.target
